# products/Telaverge/config/sut/UPF/run_shell_script_service.service
[unit]
Description = Service file
After=network.target

[Service]
Type=simple
User=root
ExecStart={{ run_smf_script }}
Restart=on-failure

[Install]
WantedBy=multi-user.target