# product/Diametriq/config/services/type1/diamsgstats.service
#/usr/lib/systemd/system/

[Unit]
Description=Message Statistics Service
After=network.target

[Service]
Type=simple
User=root
ExecStart=/bin/python /opt/remoteregal/statsscripts/diamessage_stats.py
Restart=on-failure

[Install]
WantedBy=multi-user.target

